# deploy/raspi/argos.service.example
[Unit]
Description=ARGOS V1 Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/argos/software
ExecStart=/opt/argos/software/.venv/bin/python -m argos_app
Restart=always
RestartSec=2
User=argos
Environment=ARGOS_CONFIG=/opt/argos/software/config/argos.yaml

[Install]
WantedBy=multi-user.target
